# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=/bin/bash -c "pg_isready -d db_dev -h postgres -p 5432 -U postgres && uwsgi --ini app/uwsgi.ini"
startsecs=0
autorestart=false